# source/rootfs/scripts/hica_build/chroot/loader_build/debug/udev/rules.d/11-usb-hotplug.rules
#
# Create By Czy
#

KERNEL=="sd*", SUBSYSTEM=="block", RUN+="/etc/udev/disk-hotplug.sh"
