# source/rootfs/scripts/hica_build/chroot/loader_build/release/inittab
::sysinit:/etc/init.d/rcS

::restart:/sbin/init

::shutdown:/bin/umount -a -r
