# source/rootfs/scripts/hica_build/lxc/system_build/debug/inittab
::sysinit:/etc/init.d/rcS

#::respawn:/sbin/getty -L ttyS000 115200 vt100 
::respawn:-/bin/sh

::restart:/sbin/init

#::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
